# provider/aws/sia-ec2/build/service/sia.service
[Unit]
Description=Athenz SIA
After=sshd.service

[Service]
ExecStart=/usr/sbin/siad
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
