# provider/gcp/sia-gce/build/service/sia.service
[Unit]
Description=Athenz SIA for GCE

[Service]
ExecStart=/usr/sbin/siad
Restart=always
RestartSec=120

[Install]
WantedBy=multi-user.target
