# provider/aws/sia-fargate/debian/ubuntu14/postinst
/usr/sbin/update-rc.d sia defaults
/usr/sbin/update-rc.d sia enable
