# aws-setup/ui-setup/build/service/ui.service
[Unit]
Description=Athenz UI
Wants=network-online.target ntpd.service
After=network-online.target

[Service]
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=athenz-ui
Restart=on-failure
RestartSec=20
ExecStart=/opt/athenz-ui/bin/ui.sh

[Install]
WantedBy=multi-user.target